# crontab
* * * * * curl https://stingray-app-n2l7i.ondigitalocean.app/auto-search

# An empty line is required at the end of this file for a valid cron file.
